# Blash/overlay/etc/s6-overlay/s6-rc.d/clash/run
#!/command/with-contenv sh

echo "Starting Clash ..." 

exec \
/usr/local/bin/clash -d /etc/clash 
